# tools/hardinfo2
#!/bin/sh
# Hardinfo2   support for root access
# chkconfig: - 99 10
# description: hardinfo2 support for root access

### BEGIN INIT INFO
# Provides: hardinfo2
# Required-Start: $local_fs
# Required-Stop: $local_fs
# Default-start: 2 3 4 5
# Default-Stop:
# Short-Description: start and stop hardinfo2
# Description: hardinfo2 support for root access
### END INIT INFO


start() {
   /usr/bin/hwinfo2_fetch_sysdata
   return 0
}

stop() {
   return 0
}

# See how we were called.
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  status)
	status $prog
	;;
  restart|force-reload)
	stop
	start
	;;
  try-restart|condrestart)
	if status $prog > /dev/null; then
	    stop
	    start
	fi
	;;
  reload)
	exit 3
	;;
  *)
	echo $"Usage: $0 {start|stop|status|restart|try-restart|force-reload}"
	exit 2
esac

exit 0
